# etc/supervisor/conf.d/odin_server.conf
[program:odin_server]
command = odin_server --config=%(here)s/../../odin/testing/config/test.cfg
    --static_path=%(here)s/../../odin/testing/static --debug_mode=1 --logging=info
stopsignal = INT
